# packages/addons/tools/jre/profile.d/jre.profile
export JRE_HOME=/storage/.kodi/addons/tools.jre
